# scripts/com.nuke.byok-image-analysis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nuke.byok-image-analysis</string>

    <!-- Runs in the login shell (NOT inside Claude Code's Bash sandbox), so the
         network steps reach Supabase normally. One bounded batch per fire.
         byok-image-drain.sh SELF-DRIVES across ALL the user's vehicles via a
         persistent cursor (was pinned to one hardcoded vehicle, leaving the rest
         un-analyzed); a network blip never advances the cursor or marks the fleet drained. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/skylar/nuke/scripts/daily-receipt/byok-image-drain.sh</string>
        <!-- user-id (the owner whose vehicles get drained) -->
        <string>0b9f107a-d124-49de-9ded-94698f63c1c4</string>
        <!-- batch size per fire (15 keeps an Opus-4.8 vision fire inside the 15-min window) -->
        <string>15</string>
    </array>

    <!-- 15 min between fires: paces Opus vision cost and respects the >=5min cron rule. -->
    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>/Users/skylar/nuke/logs/byok-image-batch.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/skylar/nuke/logs/byok-image-batch.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/skylar/.local/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
