# scripts/com.nuke.photo-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nuke.photo-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/skylar/nuke &amp;&amp; /usr/local/bin/dotenvx run -- /usr/local/bin/node scripts/photo-sync.mjs 2>&amp;1 | tee -a /tmp/nuke-photo-sync.log</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/tmp/nuke-photo-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/nuke-photo-sync.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
