# scripts/com.nuke.imessage-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nuke.imessage-bridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>export PATH="/opt/homebrew/bin:$PATH"; cd /Users/skylar/nuke &amp;&amp; /opt/homebrew/bin/dotenvx run -- /opt/homebrew/bin/node scripts/imessage-bridge.mjs --daemon --chat "+17026246793"</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/skylar/nuke</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/imessage-bridge.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/imessage-bridge.log</string>
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
